# streaming_video_server/deploy/systemd/streaming-video-server.user.service
[Unit]
Description=Streaming Video Server (User)
After=default.target

[Service]
Type=simple
WorkingDirectory=%h/realtime-avatar-stack/streaming_video_server
EnvironmentFile=%h/realtime-avatar-stack/streaming_video_server/.env.server
ExecStart=%h/realtime-avatar-stack/streaming_video_server/scripts/run_server.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
