# streaming_video_server/deploy/systemd/streaming-video-server.service
[Unit]
Description=Streaming Video Server
After=network.target

[Service]
Type=simple
User=streaming-video
WorkingDirectory=%h/realtime-avatar-stack/streaming_video_server
EnvironmentFile=%h/realtime-avatar-stack/streaming_video_server/.env.server
ExecStart=%h/realtime-avatar-stack/streaming_video_server/scripts/run_server.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
